# src/lib/roles/loadbalance/templates/etc/init/nginx.conf
description "nginx"
start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [016]
expect fork
kill timeout 16
respawn
respawn limit 4 16
limit nofile 16384 16384

env PID=/var/run/nginx.pid

pre-start script
  /opt/bin/nginx -t || exit $?
end script

exec /opt/bin/nginx
