# src/lib/roles/message/templates/etc/init/kafka.conf
description "kafka"

{% if inventory_hostname == groups.message[0] %}
# This is also a zookeeper host, depend upon the service appropriately.
start on started zookeeper
stop on stopping zookeeper
{% else %}
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
{% endif %}
kill timeout 16
respawn
respawn limit 2 16
limit nofile 16384 16384

pre-start script
  if [ ! -f /var/opt/kafka/log/kafka.log ]; then
    (cd /var/opt/kafka/log; touch kafka.log && chown kafka:kafka kafka.log && chmod 750 kafka.log)
  fi

  # TODO: Remove sleep once https://github.com/cloud-elements/hykes/issues/64 has been resolved.
  {% if inventory_hostname == groups.message[0] %}
  awaitd 127.0.0.1:2181 32 && sleep 4
  {% else %}
  awaitd {{ hostvars[groups.message[0]].ansible_eth0.ipv4.address}}:2181 32 && sleep 4
  {% endif %}
end script

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/kafka.pid --chuid kafka:kafka --exec /opt/kafka/bin/kafka-server-start.sh -- /etc/opt/kafka/kafka.conf >> /var/opt/kafka/log/kafka.log 2>&1
