# src/lib/roles/message/templates/etc/init/zookeeper.conf
description "zookeeper"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
kill timeout 16
respawn
respawn limit 2 16
limit nofile 16384 16384

pre-start script
  if [ ! -f /var/opt/kafka/log/zookeeper.log ]; then
    (cd /var/opt/kafka/log; touch zookeeper.log && chown kafka:kafka zookeeper.log && chmod 750 zookeeper.log)
  fi
end script

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/zookeeper.pid --chuid kafka:kafka --exec /opt/kafka/bin/zookeeper-server-start.sh -- /etc/opt/kafka/zookeeper.conf >> /var/opt/kafka/log/zookeeper.log 2>&1
