# src/lib/roles/cache/templates/etc/init/memcached.conf
description "memcached"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
respawn
respawn limit 4 16
limit nofile 16384 16384

pre-start script
  if [ ! -f /var/opt/memcached/log/memcached.log ]; then
    (cd /var/opt/memcached/log; touch memcached.log && chown memcached:memcached memcached.log && chmod 750 memcached.log)
  fi
end script

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/memcached.pid --chuid memcached:memcached --exec /opt/bin/memcached -- -u memcached -m 128 >> /var/opt/memcached/log/memcached.log 2>&1
