# src/lib/roles/cache/templates/etc/init/redis.conf
description "redis"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
expect fork
respawn
respawn limit 4 16
limit nofile 16384 16384

pre-start script
  if [ -d /sys/kernel/mm/transparent_hugepage ] && [ ! -f /sys/kernel/mm/transparent_hugepage/enabled ]; then
    echo never > /sys/kernel/mm/transparent_hugepage/enabled
  fi

  if [ ! -f /var/opt/redis/log/redis.log ]; then
    (cd /var/opt/redis/log; touch redis.log && chown redis:redis redis.log && chmod 750 redis.log)
  fi
end script

exec start-stop-daemon --start --chuid redis:redis --exec /opt/bin/redis-server -- /etc/opt/redis/redis.conf
