# src/lib/roles/logstore/templates/etc/init/elasticsearch.conf
description "elasticsearch"
start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [016]
respawn
respawn limit 4 16
limit nofile 16384 16384

{% if (((ansible_memtotal_mb / (roles_count | int)) - 512) | round(0, 'floor') | int) > 31000 %}
env ES_HEAP_SIZE=31000m
{% else %}
env ES_HEAP_SIZE={{ ((ansible_memtotal_mb / (roles_count | int)) - 512) | round(0, 'floor') | int }}m
{% endif %}

pre-start script
  if [ ! -f /var/opt/elasticsearch/log/elasticsearch.log ]; then
    (cd /var/opt/elasticsearch/log; touch elasticsearch.log && chown elasticsearch:elasticsearch elasticsearch.log && chmod 750 elasticsearch.log)
  fi
end script

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/elasticsearch.pid --chuid elasticsearch:elasticsearch --exec /opt/elasticsearch/bin/elasticsearch -- -Des.default.path.home=/opt/elasticsearch -Des.default.path.conf=/etc/opt/elasticsearch >> /var/opt/elasticsearch/log/elasticsearch.log 2>&1
