# src/lib/roles/logpipe/templates/etc/init/logstash.conf
description "logstash"
start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [016]
respawn
respawn limit 4 16
limit nofile 16384 16384

{% if (((ansible_memtotal_mb / (roles_count | int)) - 512) | round(0, 'floor') | int) > 31000 %}
env LS_HEAP_SIZE=31000m
{% else %}
env LS_HEAP_SIZE={{ ((ansible_memtotal_mb / (roles_count | int)) - 512) | round(0, 'floor') | int }}m
{% endif %}

pre-start script
  if [ ! -f /var/opt/logstash/log/logstash.log ]; then
    (cd /var/opt/logstash/log; touch logstash.log && chown logstash:logstash logstash.log && chmod 750 logstash.log)
  fi
end script

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/logstash.pid --chuid logstash:logstash --exec /opt/logstash/bin/logstash -- -f /etc/opt/logstash/logstash.conf >> /var/opt/logstash/log/logstash.log 2>&1
