# src/lib/roles/app/templates/etc/init/tomcat.conf
description "tomcat"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
kill timeout 16
respawn
respawn limit 2 16
limit nofile 16384 16384
instance $webapp
env webapp=elements

pre-start script
  if [ ! -f /var/opt/tomcat/log/tomcat.log ] ; then
    (cd /var/opt/tomcat/log ; touch tomcat.log && chown tomcat:tomcat tomcat.log && chmod 750 tomcat.log)
  fi
end script

script
  export CATALINA_HOME=/opt/tomcat
  export CATALINA_BASE=/srv/opt/tomcat/$webapp
  export CATALINA_OPTS="-Dcom.sun.management.jmxremote -Dsoba.system.dir=/etc/opt/tomcat/$webapp -Dsoba.filename=$webapp.properties";
  export JAVA_OPTS="-Xms{{ ((ansible_memtotal_mb / (roles_count | int)) * 0.55) | int }}m -Xmx{{ ((ansible_memtotal_mb / (roles_count | int)) * 0.55) | int }}m -XX:ErrorFile=/var/opt/tomcat/log/$webapp/%p.error.log -XX:+HeapDumpOnOutOfMemoryError -XX:HeapDumpPath=/var/opt/tomcat/log/$webapp/%p.heap.log -Djava.awt.headless=true -Djava.security.egd=file:/dev/./urandom -XX:+UseParNewGC -XX:+UseConcMarkSweepGC -XX:+CMSClassUnloadingEnabled -XX:ParallelGCThreads=2 -XX:NewSize={{ ((ansible_memtotal_mb / (roles_count | int)) * 0.18) | int }}m -XX:MaxNewSize={{ ((ansible_memtotal_mb / (roles_count | int)) * 0.18) | int }}m -Dfile.encoding=UTF-8 -Djsse.enableSNIExtension=false"

  rm -rf /srv/opt/tomcat/$webapp/conf
  rm -rf /srv/opt/tomcat/$webapp/logs
  rm -rf /srv/opt/tomcat/$webapp/temp
  rm -rf /srv/opt/tomcat/$webapp/work
  ln -s /etc/opt/tomcat/$webapp /srv/opt/tomcat/$webapp/conf
  ln -s /var/opt/tomcat/log/$webapp /srv/opt/tomcat/$webapp/logs
  ln -s /var/opt/tomcat/tmp/$webapp /srv/opt/tomcat/$webapp/temp
  ln -s /var/opt/tomcat/work/$webapp /srv/opt/tomcat/$webapp/work

  exec start-stop-daemon --start --make-pidfile --pidfile /var/run/tomcat.pid --chuid tomcat:tomcat --startas /bin/bash -- -c "/opt/tomcat/bin/catalina.sh run >> /var/opt/tomcat/log/tomcat.log 2>&1"
end script
